# backup.service
[Unit]
Description=Backup

[Service]
Type=oneshot
ExecStart=%h/src/backup/backup.sh
WorkingDirectory=%h/src/backup
